# MT7630E/mt7630e.service
[Unit]
Description=MT7630E

[Service]
ExecStart=/usr/local/MT7630E/load.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
